# l1pro/pwm-fan.service
[Unit]
Description=PWM fan auto speed daemon
After=network.target

[Service]
ExecStart=/usr/bin/pwm-fan.pl

[Install]
WantedBy=multi-user.target
